# .bst/staged-junctions/kde-build-meta.bst/9dc8f425e0c88ea87b3d45a0f31a74d7fc0b17758a4b174703ae285639f7ef0e/elements/kde/config/create-user.service
[Unit]
Description=Create user
ConditionFirstBoot=yes

[Service]
Type=oneshot
ExecStart=useradd -m konqi
ExecStart=usermod -a -G wheel -G video konqi

[Install]
WantedBy=multi-user.target
